# roles/dynamic-route/templates/dynamic-route.service.j2
[Unit]
Description=Dynamic Route Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/dynamic-route-parallel.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
